# docker/etc/s6/ddns/run
#!/bin/bash

# 配置
source ./setup

# 启动机器
su-exec "${USERNAME}" /usr/bin/ddns --conf "${USER_HOME}"/application.yml serve
