# docker/etc/s6/ddns/run
#!/bin/bash

# 配置
source ./setup

# 启动Archtech服务器
exec gosu "${USERNAME}" /opt/itcoursee/server --conf "${USER_HOME}"/application.yml serve
