# docker/etc/s6/ddns/run
#!/bin/bash

# 配置
source ./setup

# 启动机器
su-exec "${USERNAME}" cat "${USER_HOME}/application.yml"
su-exec "${USERNAME}" /opt/storezhang/ddns --conf "${USER_HOME}/application.yml" serve
