# docker/etc/s6/ddns/run
#!/bin/bash

# 配置
source ./setup

# 启动机器
su-exec "${USERNAME}" /opt/storezhang/ddns -c="${USER_HOME}/application.yml" serve
