# docker/etc/s6/songjiang/run
#!/bin/bash

[[ -f ./setup ]] && source ./setup

exec su-exec "${USERNAME}" /opt/songjiang --conf /conf/application.yml
