# root/etc/s6/ddns/run
#!/bin/bash

[[ -f ./setup ]] && source ./setup

exec su-exec "${USERNAME}" /opt/ddns --conf /conf/ddns.yml
